# network/services/quickstart/files/startdhns.init
#!/bin/sh /etc/rc.common

START=91
USE_PROCD=1

boot() {
	start "BOOT"
}

start_service() {
	[ "$1" = BOOT ] && return 0
	/usr/sbin/quickstart uciChange
}

service_triggers()
{
	procd_add_reload_trigger network
}
